# bin/ai-newsletter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pai.ai-newsletter</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/robt/.bun/bin/bun</string>
        <string>/Users/robt/.claude/bin/ai-newsletter</string>
        <string>--publish</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <!-- Run daily at 7:00 AM -->
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/robt/.claude/logs/ai-newsletter.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/robt/.claude/logs/ai-newsletter.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/robt/.bun/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/robt</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
